# ops/systemd/shuanglu-certbot-renew.service
[Unit]
Description=Renew Shuanglu Let's Encrypt certificate
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker run --rm --network host --volume /etc/letsencrypt:/etc/letsencrypt --volume /home/web/letsencrypt:/var/www/letsencrypt certbot/certbot:v5.6.0 renew --quiet --no-random-sleep-on-renew --cert-name shuanglu.uway.click
ExecStartPost=/usr/local/sbin/deploy-shuanglu-certificate
