# ops/systemd/shuanglu.service
[Unit]
Description=Shuanglu Next.js application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=shuanglu
Group=shuanglu
WorkingDirectory=/opt/apps/shuanglu/current
Environment=NODE_ENV=production
Environment=HOME=/opt/apps/shuanglu
ExecStart=/usr/bin/npm start -- --hostname 127.0.0.1 --port 3002
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
MemoryMax=1G
CPUQuota=100%
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
